# etc/crontabs/root
33 * * * * logrotate /etc/logrotate.d/aprx > /tmp/crontab
